# configs/systemd/kioskbook-app.service
[Unit]
Description=KioskBook Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kioskbook
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=5173
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
