# configs/systemd/kioskbook-app.service
[Unit]
Description=KioskBook Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=kioskbook-app
Group=kioskbook-app
WorkingDirectory=/opt/kioskbook/dist
Environment=NODE_ENV=production
ExecStart=/usr/bin/npx serve -s . -l tcp://127.0.0.1:5173
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/kioskbook

[Install]
WantedBy=multi-user.target
